# flantier.service
[Unit]
Description=Flantier der Geschenk Manager
After=network-online.target

[Service]
Type=simple
ExecStart=/home/ubuntu/.local/bin/flantier

User=ubuntu
Group=ubuntu

Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
