# flantier.service
[Unit]
Description=Flantier der Geschenk Manager
After=network-online.target

[Service]
Type=simple
ExecStart=flantier

Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
